# fleet-worker/artifacts-v1/ups-epd.service
[Unit]
Description=UPS E-Paper Runtime Display
After=nut-server.service
Requires=nut-server.service

[Service]
Type=simple
Environment=UPS_NAME=ups
Environment=INTERVAL_AC=10
Environment=INTERVAL_BATT=5
Environment=FULL_REFRESH_EVERY=0
Environment=DISPLAY_TZ=America/Los_Angeles
Environment=TRANSITION_DELAY_TO_BATT=1.0
Environment=TRANSITION_DELAY_TO_AC=10.0
Environment=ESTIMATE_ONLINE_RUNTIME=1
Environment=MODEL_STATE_PATH=/var/lib/ups-epd/runtime_model.json
Environment=MODEL_MIN_WATTS=20
Environment=UPS_RATED_WATTS=0
Environment=BOGUS_OL_RUNTIME_THRESHOLD_S=300
ExecStart=/usr/bin/python3 /opt/ups-epd/ups_epd.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
